# Configs/etc/rc.local
# NOTE: not part of newer distributions see /etc/NetworkManager/dispatcher.d

source /etc/iptables.rules
