# Configs/home/.config/autostart/screen_lock.desktop
[Desktop Entry]
Type=Application
Name=Screen Lock
Exec=gnome-screensaver-comand -l
